# post_script_payload/to_install/etc/systemd/system/set-up-and-update.service
[Unit]
Description=Set up and then update system
Wants=network-online.target
After=network-online.target
FailureAction=reboot
SuccessAction=reboot

[Service]
ExecStart=/usr/local/bin/set-up-and-update.sh
