# post_script_payload/to_install/etc/systemd/system/set-up-and-update.target
[Unit]
Description=Minimal environment to set up and update system
Wants=set-up-and-update.service
