# updates-phase-1.target
[Unit]
Description=Activate phase 1 of system updates on boot
Wants=updates-phase-1.service
